# 12_Automatisering_med_Configuration_Management_program/VG_ansible/roles/web/files/demod.service
[Unit]
Description=Ansible demo app.

[Service]
ExecStart=/var/www/demo/bin/app.pl --port=8080 --environment=production

[Install]
WantedBy=multi-user.target